# com.a58.okproperty.dashboard.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.a58.okproperty.dashboard.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/a58/claude专用/okproperty_dashboard/run_sync.sh</string>
        <string>30</string>
        <string>auto</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/a58/claude专用/okproperty_dashboard</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.10/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- 每天 08:00 主跑；09:00 / 10:00 兜底(auto 模式：当天已成功则跳过)。10点也失败→飞书回执 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/a58/claude专用/okproperty_dashboard/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/a58/claude专用/okproperty_dashboard/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
